# debian/opi-postsetup.cron.d
#
# Run firstboot scripts only at reboot
#
@reboot root /usr/share/opi-postsetup/kgp-postsetup
